# scripts/genie-rs.service
[Unit]
Description=Genie-RS OpenAI-compatible API for Dragon NPU
After=network.target

[Service]
Type=simple
User=${USER}
Environment=QAIRT=${QAIRT_SDK:-/opt/qairt}
Environment=LD_LIBRARY_PATH=${MODEL_DIR}
Environment=PATH=${CARGO_HOME:-$HOME/.cargo}/bin:/usr/bin:/bin
ExecStart=${GENIE_RS_BIN:-/opt/genie-rs}/genie-rs serve \
    --host 0.0.0.0:8080 \
    --registry ${REGISTRY:-/etc/genie-rs/registry.toml}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
